# lumina-wallpaper-preseed.service
[Unit]
Description=Pre-populate SDDM lumina wallpaper before login screen
Before=sddm.service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'WALL=$(cat /home/cuong/.local/state/caelestia/wallpaper/path.txt 2>/dev/null); [ -n "$WALL" ] && [ -f "$WALL" ] && cp "$WALL" /var/cache/lumina-sddm/wallpaper.jpg || true'
RemainAfterExit=yes

[Install]
WantedBy=sddm.service
