# lumina-wallpaper-preseed.service
[Unit]
Description=Pre-populate SDDM lumina wallpaper before login screen
Before=sddm.service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'MAINUSER=$(getent passwd | awk -F: '"'"'$3>=1000 && $3<65000{print $1;exit}'"'"'); STATE=/home/$MAINUSER/.local/state; WALL=$(cat $STATE/caelestia/wallpaper/path.txt 2>/dev/null); [ -n "$WALL" ] && [ -f "$WALL" ] && cp "$WALL" /var/cache/lumina-sddm/wallpaper.jpg || true'
RemainAfterExit=yes

[Install]
WantedBy=sddm.service
